# services/set-sleep/auto-sleep.timer
[Unit]
Description=Periodic check for idle state (auto-sleep)

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
